# engine/install-ovirt.service
[Unit]
Description=Install oVirt

[Service]
ExecStart=/install-ovirt.sh

[Install]
WantedBy=multi-user.target
